# debian/nodejs-pool-block-manager.service
[Unit]
Description=Mining poor for Monero - Block Manager
After=network.target

[Service]
WorkingDirectory=/usr/share/nodejs-pool
ExecStart=/usr/bin/nodejs /usr/share/nodejs-pool/init.js --name=blockManager --log-date-format="YYYY-MM-DD HH:mm Z" --module=blockManager
Restart=always
User=nodejs-pool

[Install]
WantedBy=multi-user.target
