# debian/nodejs-pool-long-runner.service
[Unit]
Description=Mining poor for Monero - Long Runner
After=network.target

[Service]
WorkingDirectory=/usr/share/nodejs-pool
ExecStart=/usr/bin/nodejs /usr/share/nodejs-pool/init.js --name=longRunner --log-date-format="YYYY-MM-DD HH:mm Z" --module=longRunner
Restart=always
User=nodejs-pool

[Install]
WantedBy=multi-user.target
