# debian/nodejs-pool-api.service
[Unit]
Description=Mining poor for Monero - API
After=network.target

[Service]
WorkingDirectory=/usr/share/nodejs-pool
ExecStart=/usr/bin/nodejs /usr/share/nodejs-pool/init.js --name=api --log-date-format="YYYY-MM-DD HH:mm Z" --module=api
Restart=always
User=nodejs-pool

[Install]
WantedBy=multi-user.target
