# debian/nodejs-pool-pool.service
[Unit]
Description=Mining poor for Monero - Pool
After=network.target

[Service]
WorkingDirectory=/usr/share/nodejs-pool
ExecStart=/usr/bin/nodejs /usr/share/nodejs-pool/init.js --name=pool --log-date-format="YYYY-MM-DD HH:mm Z" --module=pool
Restart=always
User=nodejs-pool

[Install]
WantedBy=multi-user.target
